# systemd/opencode-telegram-bridge.service
[Unit]
Description=OpenCode Telegram Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/opencode-telegram-bridge/opencode-telegram-bridge.env
ExecStart=/usr/bin/opencode-telegram-bridge
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
